# systemd/airquality-poller.timer
[Unit]
Description=Run Air Quality Poller every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
