# systemd/airquality.service
[Unit]
Description=Air Quality Streamlit Dashboard
After=network.target

[Service]
Type=simple
User=exouser
WorkingDirectory=/home/exouser/air-quality-report
EnvironmentFile=/home/exouser/air-quality-report/.env
ExecStart=/home/exouser/air-quality-report/venv/bin/streamlit run dashboard.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
