# systemd/airquality-poller.service
[Unit]
Description=Air Quality Poller (single run)
After=network.target

[Service]
Type=oneshot
User=exouser
WorkingDirectory=/home/exouser/air-quality-report
EnvironmentFile=/home/exouser/air-quality-report/.env
ExecStart=/home/exouser/air-quality-report/venv/bin/python poller.py
StandardOutput=journal
StandardError=journal
